# games-action/mt_test/files/minetestserver.initd
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Copyright 2019 Wade Cline
# Distributed under the terms of the GNU General Public License v2

description="Minetest dedicated server"

# Restart on crash.
command="${MINETESTBIN}"
command_args="${ARGS}"
command_user="${USER}:${GROUP}"
pidfile="${PIDFILE}"
supervisor=supervise-daemon

depend() {
	need net
}

