# app-backup/minetest-backup/files/mtbackup.cron
PATH="/sbin:/bin:/usr/sbin:/usr/bin"
10 3 * * * mtbackup
